# finance-ai.service
[Unit]
Description=Finance AI backend service
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/workspace/finance-ai
WorkingDirectory=/home/devops/finance-ai
RemainAfterExit=yes
ExecStart=/usr/bin/docker compose -f /workspace/finance-ai/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f /workspace/finance-ai/docker-compose.yml down
ExecStart=/usr/bin/docker compose -f /home/devops/finance-ai/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f /home/devops/finance-ai/docker-compose.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target